# configs/host_agnostic/.local/etc/anacrontab
@monthly ${HOME}/bin/new_accounting_month
